# usr/etc/inittab
chmod 04755 /bin/passwd
chmod 04755 /bin/su
/bin/getty
